# x-feed-intel-dashboard.service
[Unit]
Description=X Feed Intel Dashboard
After=network.target

[Service]
Type=simple
User=xfi
Group=xfi
WorkingDirectory=/opt/x-feed-intel
Environment="PATH=/opt/x-feed-intel/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="XFI_COMMON_CTX_STATS_TTL_SEC=5"
Environment="XFI_SLOW_REQUEST_LOG_MS=350"
ExecStart=/opt/x-feed-intel/venv/bin/waitress-serve --listen=0.0.0.0:5050 --threads=2 dashboard:app
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Resource limits (adjust for your deployment)
MemoryMax=256M
CPUQuota=50%
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
