# x-feed-intel-dashboard.service
[Unit]
Description=X Feed Intel Dashboard
After=network.target

[Service]
Type=simple
User=__RUN_AS_USER__
Group=__RUN_AS_GROUP__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=-__ENV_FILE__
Environment="XFI_COMMON_CTX_STATS_TTL_SEC=5"
Environment="XFI_SLOW_REQUEST_LOG_MS=350"
ExecStart=__VENV_BIN_DIR__/waitress-serve --listen=0.0.0.0:5050 --threads=2 dashboard:app
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Resource limits (adjust for your deployment)
MemoryMax=256M
CPUQuota=50%
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
